# chromeos-base/chromeos-bsp-surface-pro4/files/tablet-mode/99-usb-keyboard-detect.rules
ACTION=="add", SUBSYSTEMS=="usb", ATTRS{idVendor}=="045e", ATTRS{idProduct}=="09c0", RUN+="switch_tablet_mode.sh 0"
ACTION=="remove", SUBSYSTEMS=="usb", ATTRS{idVendor}=="045e", ATTRS{idProduct}=="09c0", RUN+="switch_tablet_mode.sh 1"
